# BirthdayBox/crontab
#run python script every minutes
*/15 * * * * cd /usr/src/app; bash run.sh > /proc/1/fd/1 2>/proc/1/fd/2
